# deploy/systemd/lng-cleanup-local-raw.timer
[Unit]
Description=Run lng-cleanup-local-raw.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
